# .archive/aws-sam-cli/services.d/10-dynamodb-admin/run
#!/usr/bin/with-contenv bash

echo "Starting dynamodb-admin ..."
exec \
    dynamodb-admin

echo "Started - dynamodb-admin ui"